# scripts/supervisord.conf
[supervisord]
nodaemon=true
user=postal
logfile=/dev/stdout
logfile_maxbytes=0
pidfile = /tmp/supervisord.pid

[inet_http_server]
port = 0.0.0.0:9001
username = %(ENV_SD_USER)s
password = %(ENV_SD_PASS)s

[program:postal]
command=/opt/postal/bin/postal run
user=postal
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
